# config/cloudquery/cloudquery.service
[Unit]
Description=CloudQuery

[Service]
ExecStart=/cloudquery sync /aws.yaml /postgresql.yaml

[Install]
WantedBy=multi-user.target