# config/cloudquery/cloudquery.service
[Unit]
Description=CloudQuery

[Service]
ExecStart=/opt/cloudquery/cloudquery.sh

[Install]
WantedBy=multi-user.target
